# mqtt-sensor-daemon.service
[Unit]
Description=MQTT Sensor Daemon
After=network-online.target
Wants=network-online.target
Requires=pigpiod.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/mqtt-sensor-daemon
ExecStart=/opt/mqtt-sensor-daemon/venv/bin/python /opt/mqtt-sensor-daemon/mqtt-sensor-daemon.py /opt/mqtt-sensor-daemon/config.ini
Environment=PYTHONUNBUFFERED=1
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target